# scripts/crontab
* * * * * echo "I'm alive" >> /var/log/cron.log 2>&1

30 07 * * * /root/code/scripts/run.sh >> /var/log/cron.log 2>&1

#Don't remove new line




